# backend/deployment/das-tool-rest.service
[Unit]
Description=Das Tool Rest API
After=syslog.target

[Service]
User=oblume
ExecStart=/home/oblume/das-tool-rest/das-tool-rest.jar
WorkingDirectory=/home/oblume/das-tool-rest
SuccessExitStatus=143

[Install]
WantedBy=multi-user.target
